# voice.service
[Unit]
Description=Smart Home Voice Intelligence
After=network-online.target hub.service
Wants=network-online.target

[Service]
Type=simple
User=anadivatsa
WorkingDirectory=/home/anadivatsa/smarthome
ExecStart=/home/anadivatsa/smarthome/venv/bin/python voice.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/home/anadivatsa/smarthome/voice.env
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
